# ansible/templates/traefik/traefik.service.j2
[Unit]
Description=Traefik Reverse Proxy
Documentation=https://doc.traefik.io/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=traefik
Group=traefik

# Environment variables for ACME/Cloudflare
Environment="CLOUDFLARE_DNS_API_TOKEN={{ cloudflare_dns_token }}"
Environment="CLOUDFLARE_ZONE_API_TOKEN={{ cloudflare_dns_token }}"

# Traefik executable
ExecStart=/usr/local/bin/traefik --configFile=/etc/traefik/traefik.toml

# Process management
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=traefik

# Security - minimal sandboxing to avoid startup issues
NoNewPrivileges=true

# Resource limits (adjust as needed)
LimitNOFILE=1000000
LimitNPROC=512

[Install]
WantedBy=multi-user.target
